# dness-bin/dness.timer
[Unit]
Description=Pacman Auto Update Timer

[Timer]
OnBootSec=5min
OnUnitInactiveSec=5min
Unit=dness.service

[Install]
WantedBy=multi-user.target
