# dmarc-report-converter/dmarc-report-converter.service
[Unit]
Description=DMARC report converter
After=network.target

[Service]
Type=simple
User=nginx
Group=nginx
ExecStart=/usr/bin/dmarc-report-converter
KillMode=process

[Install]
WantedBy=multi-user.target
