# audirvana-studio/audirvanaStudio.service
[Unit]
Description=Run audirvanaStudio
DefaultDependencies=no
After=avahi-daemon.service

[Service]
Type=simple
WorkingDirectory=/opt/audirvana/studio
ExecStart=/opt/audirvana/studio/audirvanaStudio
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
