# hoarder/hoarder-browser.service
[Unit]
Description = Hoarder browser
Wants = network-online.target
After = network-online.target

[Service]
User = hoarder
Restart = on-failure
ExecStart = /usr/bin/chromium --headless --no-sandbox --disable-gpu --disable-dev-shm-usage --remote-debugging-address=127.0.0.1 --remote-debugging-port=9222 --hide-scrollbars
TimeoutStopSec = 5

[Install]
WantedBy = multi-user.target
