# auto_facelock/auto_facelock@.service
[Unit]
Description=Face Detection Screen Lock Service for %I
After=graphical.target
Requires=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto_facelock.py --known_faces_dir /usr/local/share/known_faces --check_interval 60 --retries 2
Restart=always
RestartSec=60
User=%I
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%I/.Xauthority

[Install]
WantedBy=default.target
