# ceserver/ceserver.service
[Unit]
Description=Cheat Engine Server
After=network.target

[Service]
ExecStart=/opt/ceserver/ceserver
Restart=on-failure

[Install]
WantedBy=multi-user.target
