# ollama-copilot-git/ollama-copilot.service
[Unit]
Description=Start the ollama-copilot server
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/ollama-copilot

[Install]
WantedBy=multi-user.target
