# python-openlp-osc-control/OpenLPOSCControl.service
[Unit]
Description=OpenLP OSC Control Service

[Service]
ExecStart=/usr/bin/openLPOSCControl

[Install]
WantedBy=default.target
