# deepcool-digital-linux-git/deepcool-digital.service
[Unit]
Description=DeepCool Digital

[Service]
ExecStart=/usr/sbin/deepcool-digital-linux

[Install]
WantedBy=multi-user.target