# paperless-ngx-venv/paperless-webserver.service
[Unit]
Description=Paperless webserver service
After=network.target
Requires=redis.service
PartOf=paperless.target

[Service]
User=paperless
Group=paperless
EnvironmentFile=/etc/paperless.conf
WorkingDirectory=/usr/share/paperless/src
ExecStart=/usr/lib/paperless/webserver
Restart=on-abort

[Install]
WantedBy=paperless.target
