# audirvana-origin/audirvanaOrigin.service
[Unit]
Description=Run audirvanaOrigin
DefaultDependencies=no
After=avahi-daemon.service

[Service]
Type=simple
WorkingDirectory=/opt/audirvana/origin
ExecStart=/opt/audirvana/origin/audirvanaOrigin
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
