# flexget/flexget.service
[Unit]
Description=FlexGet Daemon (system)
Documentation=https://www.flexget.com
After=network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/flexget daemon start --autoreload-config
ExecStop=/usr/bin/flexget daemon stop
ExecReload=/usr/bin/flexget daemon reload
KillMode=process
User=flexget
WorkingDirectory=/var/lib/flexget

[Install]
WantedBy=multi-user.target
