# chinesesubfinder/chinesesubfinder.user.service
[Unit]
Description = Automated Chinese subtitle download.
Wants = network-online.target
After = network-online.target

[Service]
Restart = on-failure
ExecStart = /usr/bin/chinesesubfinder --setconfigselfpath %h/.local/share/chinesesubfinder

[Install]
WantedBy = default.target
