# tuned-openrc/tuned-ppd
#!/sbin/openrc-run

name=$RC_SVCNAME
description="PPD-to-TuneD API Translation Daemon"
supervisor="supervise-daemon"
command="/usr/bin/tuned-ppd"
depend() {
	after tuned 
	before multi-user display-manager 
	need tuned 
}
