# jackett-mono/jackett.service
[Unit]
Description=Jackett Daemon
After=network.target

[Service]
User=jackett
Group=jackett
SyslogIdentifier=jackett
Restart=always
RestartSec=5
Type=simple
Environment=XDG_CONFIG_HOME=/var/lib/jackett
Environment=DOTNET_EnableDiagnostics=0
WorkingDirectory=/usr/lib/jackett
ExecStart=/usr/bin/mono /usr/lib/jackett/JackettConsole.exe --NoRestart --NoUpdates --DataFolder /var/lib/jackett
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
