# borgwarehouse/borgwarehouse-cron.timer
[Unit]
Description=Borgwarehouse periodic job

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
