# odoo-nightly/odoo.service
[Unit]
Description=Odoo Open Source ERP and CRM
Requires=postgresql.service
After=postgresql.service network-online.target

[Service]
Type=simple
User=odoo
Group=odoo
ExecStart=/var/lib/odoo/.venv/bin/odoo --config /etc/odoo/odoo.conf 
KillMode=mixed

[Install]
WantedBy=multi-user.target
