# nextpvr-bin/nextpvr.service
[Unit]
Description=NextPVR video channel recording

[Install]
WantedBy=default.target

[Service]
Type=simple
ExecStart=/usr/bin/nextpvr

