# sshesame/sshesame@.service
[Unit]
Description=An easy to set up and use SSH honeypot, a fake SSH server that lets anyone in and logs their activity
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/sshesame -config /etc/sshesame/config.yaml
Restart=always

[Install]
WantedBy=multi-user.target