# fusioninventory-agent-rust-git/fusioninventory-agent.service
[Unit]
Description=FusionInventory Rust agent
After=network.target

[Service]
WorkingDirectory=/opt/fusioninventory/
ExecStart=/opt/fusioninventory/fusioninventory-agent -d

[Install]
WantedBy=multi-user.target