# elasticsearch-bin/elasticsearch-keystore.service
[Unit]
Description=Elasticsearch Keystore Generation
ConditionPathExists=|!/etc/elasticsearch/elasticsearch.keystore

[Service]
Environment=ES_PATH_CONF=/etc/elasticsearch
Type=oneshot
Group=elasticsearch
UMask=0007
ExecStart=/usr/share/elasticsearch/bin/elasticsearch-keystore create
RemainAfterExit=yes
